# deploy/systemd/pybulletin-web.service
[Unit]
Description=pyBulletin web service (public BBS UI + sysop console)
After=network-online.target pybulletin.service
Wants=network-online.target

[Service]
Type=simple
User=pybulletin
Group=pybulletin
WorkingDirectory=/home/pybulletin/pyBulletin
Environment=PYTHONPATH=src
ExecStart=/usr/local/bin/pybulletin-python -m pybulletin.cli --config config/pybulletin.toml serve-web
Restart=on-failure
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=10
SendSIGKILL=yes

[Install]
WantedBy=multi-user.target
